# clearLog.service
[Unit]
Description=Clears the chat log

[Service]
Type=oneshot
ExecStart=/usr/bin/clearLog.sh

[Install]
RequiredBy=multi-user.target
